# board/qemu/mipsel/overlay/etc/sv/atmtcp/run
#!/bin/sh

set -eu

SERVER=192.0.2.0

# wait for the remote end to become available
echo | nc $SERVER 2812 2>/dev/null || exit 0

exec atmtcp virtual connect $SERVER
