# board/qemu/mipsel/overlay/etc/sv/atmtcp/run
#!/bin/sh

set -eu

exec atmtcp virtual connect 192.0.2.0
