# fakeisp/rootfs-overlay/opt/bratwurst/runit/isc-dhcp-server/run
#!/bin/sh

set -eu

exec 2>&1
exec /usr/sbin/dhcpd -4 -f -cf /opt/bratwurst/services/isc-dhcp-server/dhcpd.conf -d
