# rootfs/etc/sv/br2684ctl/run
#!/bin/sh

set -eu

. /etc/bratwurst

case $ATM_ENCAP in
llc)		ENCAP=0;;
vcmux)		ENCAP=1;;
esac

exec /usr/sbin/br2684ctl -c $UPLINK_PORT -e $ENCAP -a $UPLINK_PORT.$ATM_VPI.$ATM_VCI
