# overlay/etc/sv/watchdog/run
#!/bin/sh

set -eu

/sbin/modprobe softdog

exec /sbin/watchdog -F /dev/watchdog
