# rootfs/etc/sv/watchdog/run
#!/bin/sh

set -eu

if [ -c /dev/watchdog ]; then
	echo /dev/watchdog does not exist
	sv once .
	exit 1
fi

exec /sbin/watchdog -F /dev/watchdog
