# fakeisp/overlay/opt/bratwurst/runit/atmtcp-whack/run
#!/bin/sh

set -eu

exec nc -l -p 2813 </dev/null 2>/dev/null
