# rootfs/etc/sv/odhcp6c/run
#!/bin/sh

set -eu

. /etc/bratwurst

[ -f /run/uplink ] && . /run/uplink || exit 1

[ ! "$PREFIX" ] && PD="-P 0"

exec 2>&1
exec /usr/sbin/odhcp6c -v -e ${PD:-} -s /usr/share/dhcpc.script $IFNAME
