# rootfs/etc/sv/odhcp6c/run
#!/bin/sh

set -eu

[ -f /run/uplink ] || exit 1

. /run/uplink

exec 2>&1
exec /usr/sbin/odhcp6c -v -e -P 0 -s /usr/share/dhcpc.script $IFNAME
