# fakeisp/rootfs-overlay/opt/bratwurst/runit/atmtcp-whack/finish
#!/bin/sh

set -eu

sv force-restart atmtcp || true

exit 0
