# fakeisp/rootfs-overlay/opt/bratwurst/runit/radvd/run
#!/bin/sh

set -eu

exec 2>&1
exec /usr/sbin/radvd -n -m stderr -C /opt/bratwurst/services/radvd/radvd.conf
