# fakeisp/rootfs-overlay/opt/bratwurst/runit/atmtcp/run
#!/bin/sh

set -eu

/sbin/modprobe atmtcp

exec 2>&1
exec /usr/sbin/atmtcp virtual listen
