# fakeisp/rootfs-overlay/opt/bratwurst/runit/pppoe/finish
#!/bin/sh

set -eu

exec sv force-stop br2684ctl
