# board/qemu/mipsel/rootfs-overlay/etc/sv/atmtcp/finish
#!/bin/sh

set -eu

SERVER=172.21.0.0

exec timeout -t 3 nc $SERVER 2813 </dev/null >/dev/null 2>/dev/null
