# rootfs/etc/sv/ntpd/run
#!/bin/sh

set -eu

[ -f /run/uplink ] || exit 1

exec 2>&1
exec /sbin/ntpd -n -l -p pool.ntp.org
