# fakeisp/overlay/opt/bratwurst/runit/atmtcp/run
#!/bin/sh

set -eu

/sbin/modprobe atmtcp

exec atmtcp virtual listen
